# deploy/systemd/retriever-web.service.example
[Unit]
Description=Retriever web service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=retriever
Group=retriever
WorkingDirectory=/opt/retriever-rebuild/current
EnvironmentFile=/etc/retriever-rebuild/retriever.env
ExecStart=/opt/retriever-rebuild/current/.venv/bin/uvicorn app.main:app --host ${RETRIEVER_BIND_HOST} --port ${RETRIEVER_PORT}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

